# trip_ingestion_service/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:ingestor]
command=python /app/ingestor.py
autostart=true
autorestart=true

[program:pruner]
command=python /app/pruner.py
autostart=true
autorestart=true
